# fs/etc/services.d/pushgateway/run
#!/usr/bin/with-contenv ash

cd /pushgateway

exec runas /bin/pushgateway $PUSHGATEWAY_ARGS